# mac/com.rock.ai-usage-push.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.rock.ai-usage-push</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/rock/code/esp32-ai-usage-display/mac/push-usage.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>60</integer>
  <key>StandardOutPath</key>
  <string>/Users/rock/Library/Logs/ai-usage-push.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/rock/Library/Logs/ai-usage-push.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/rock/.nvm/versions/node/v22.22.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
